# systemd/wetty.socket
[Unit]
Description=Wetty node app socket

[Socket]
ListenStream=9999

[Install]
WantedBy=sockets.target
