# upstart/wetty.conf
description "wetty terminal server"
author "me"
start on runlevel [2345]
exec nodejs /opt/b2g/wetty/app.js -p 9999 --sshauth publickey,password
